# meta-luv/recipes-core/telemetrics/files/telemetrics_script.service
[Unit]
Description=should telemetrics run
Before=telemd.service

[Service]
Type=oneshot
ExecStart=@SBINDIR@/telemetrics_runner

[Install]
WantedBy=multi-user.target
